# templates/default/influxdb.conf.erb
# InfluxDB upstart script

description "InfluxDB service"

start on (net-device-up
          and local-filesystems
          and runlevel [2345])

stop on runlevel [016]

respawn
respawn limit 10 30

limit nofile 32000 32000

setuid <%= node['influxdb']['user'] %>
setgid <%= node['influxdb']['group'] %>

<% if @ec2 %>
exec /opt/influxdb/influxd-ec2
<% else %>
exec /opt/influxdb/influxd -config /etc/opt/influxdb/influxdb.conf
<% end %>
